# systemd/claude-health-monitor.timer
[Unit]
Description=CTS Health Monitor (every 2 min)

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=15s

[Install]
WantedBy=timers.target
