# systemd/claude-heartbeat-tier3.service
[Unit]
Description=CTS Heartbeat Tier 3 — Daily Situation Report
After=network.target

[Service]
Type=oneshot
ExecStart=/home/johnlane/projects/claude-task-system/cts-v2/src/heartbeat/tier3.sh
Environment=HOME=/home/johnlane
Environment=PATH=/home/johnlane/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/johnlane/projects/claude-task-system
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
