# systemd/ralph-pro-monitor.timer
[Unit]
Description=Ralph Pro Overnight Monitor Timer (every 30 min)

[Timer]
OnCalendar=*:00/30
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
