# systemd/ralph-pro-fixer.service
[Unit]
Description=Ralph Pro Fixer (processes fix requests)
After=network.target

[Service]
Type=oneshot
ExecStart=/home/johnlane/ralph-pro/monitor/fixer.sh process
Environment=HOME=/home/johnlane
Environment=PATH=/home/johnlane/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/johnlane/ralph-pro
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
