# systemd/ralph-pro-digest.timer
[Unit]
Description=Ralph Pro Daily Digest (08:00 UTC)

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
